# .config/sv/emacs-daemon/run
#!/bin/sh

exec /usr/bin/emacs --fg-daemon
