# .config/sv/transmission-daemon/run
#!/bin/sh

exec /usr/bin/transmission-daemon -f
